# deploy/systemd/ohmysmsd.service
[Unit]
Description=ohmysmsd - SMS management daemon
Documentation=https://github.com/KimmyXYC/ohmysmsapp
After=network-online.target ModemManager.service
Wants=network-online.target
# ModemManager 必须在线才能发现模块；如果 MM 重启我们也跟着重启。
Requires=ModemManager.service
PartOf=ModemManager.service

[Service]
Type=simple
User=root
Group=root
# 这是一个示例模板，请按实际安装目录调整。
WorkingDirectory=/srv/ohmysmsapp
ExecStart=/srv/ohmysmsapp/bin/ohmysmsd -config /srv/ohmysmsapp/config.yaml
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# 需要 root 访问 DBus System bus 上的 ModemManager1（默认 policy 允许 root/网络管理组）
# 以及 /dev/cdc-wdm* /dev/ttyUSB* 等。如后续降权，改为专用用户 + polkit rules。

# 资源 & 安全
LimitNOFILE=65536
# 注意：不能开 NoNewPrivileges=true。我们 fork 的 lpac 子进程依赖
# LD_LIBRARY_PATH=/srv/ohmysmsapp/lib 或类似目录来加载 driver_apdu_*.so，
# 而 NoNewPrivileges 触发 ld.so 的 secure-execution 模式（AT_SECURE=1），
# 在该模式下 LD_LIBRARY_PATH 被忽略，导致 lpac 报 "No APDU driver found"。
# 如未来把 .so 装到 ld.so.cache 路径（/usr/local/lib/lpac + ldconfig），可以重新启用。
ProtectSystem=strict
ReadWritePaths=/srv/ohmysmsapp/data /srv/ohmysmsapp/web
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# 需要访问串口/WDM 设备；若 User=改成非 root 需加 SupplementaryGroups=dialout
# 根据实际硬件暴露的设备节点增删。
DeviceAllow=/dev/cdc-wdm* rw
DeviceAllow=/dev/ttyUSB* rw

[Install]
WantedBy=multi-user.target
